# docker/render/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
; Sprout backend — supervisord config for the single-container Render deployment.
;
; Runs all five Python microservices on internal localhost ports plus the Caddy
; reverse proxy on $PORT. All logs go to stdout/stderr so Render captures them.
; ─────────────────────────────────────────────────────────────────────────────
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; Enable `supervisorctl status` inside the container for debugging on Render:
;   render ssh   (then)   supervisorctl -c /etc/supervisord.conf status
[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; Common log settings for every program: stream to the container's stdout so
; `render logs` / `docker logs` shows everything, with no rotation.
[program:registry_api]
command=/app/.venv/bin/uvicorn sprout_registry.main:app --host 127.0.0.1 --port 8766
directory=/app
autostart=true
autorestart=true
startsecs=5
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:tool_executor]
command=/app/.venv/bin/uvicorn sprout_executor.main:app --host 127.0.0.1 --port 8767
directory=/app
autostart=true
autorestart=true
startsecs=5
priority=15
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:chat_backend]
command=/app/.venv/bin/uvicorn sprout_chat_backend.main:app --host 127.0.0.1 --port 8765
directory=/app
autostart=true
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:synthesis_service]
command=/app/.venv/bin/uvicorn sprout_synthesis.main:app --host 127.0.0.1 --port 8002
directory=/app
autostart=true
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:mcp_server]
command=/app/.venv/bin/python -m sprout_mcp.main streamable-http
directory=/app
environment=SPROUT_MCP_HOST="127.0.0.1",SPROUT_MCP_PORT="8768"
autostart=true
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

; Caddy binds the public $PORT and routes to the services above. Start last.
[program:caddy]
command=/usr/bin/caddy run --config /etc/caddy/Caddyfile --adapter caddyfile
directory=/app
autostart=true
autorestart=true
startsecs=3
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
